# systemd/taos-pre-shutdown.service
[Unit]
Description=taOS graceful agent shutdown hook
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/taos-graceful-stop
TimeoutStopSec=60

[Install]
WantedBy=shutdown.target reboot.target halt.target
