# systemd/tinyagentos.service
[Unit]
Description=TinyAgentOS Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/tinyagentos
ExecStart=%h/tinyagentos/.venv/bin/python -m uvicorn tinyagentos.app:create_app --factory --host 0.0.0.0 --port 6969
ExecStop=/usr/local/bin/taos-graceful-stop
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
