# os-build/userpatches/overlay/etc/systemd/system/tinyagentos.service
[Unit]
Description=TinyAgentOS Web GUI
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=taos
Group=taos
WorkingDirectory=/opt/tinyagentos
ExecStart=/opt/tinyagentos/venv/bin/python -m uvicorn tinyagentos.app:create_app --factory --host 0.0.0.0 --port 6969
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
